# racket-api.service
[Unit]
Description=Racket API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/curl-tut
ExecStart=/usr/bin/racket /var/www/curl-tut/server.rkt
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
